# vault/service/consul/run
#!/bin/sh
consul agent -config-file=/etc/consul.json
